# botsbegone.service
[Unit]
Description=BotsBeGone firewall service update

[Service]
Type=simple
ExecStart=/var/lib/BotsBeGone/BotsBeGone/botsbegone.py
Restart=no