# botsbegone.timer
[Unit]
Description=BotsBeGone firewall service update

[Timer]
OnCalendar=*-*-* 00:00:00
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target


